# cron.d/xentral
* * * * * www-data php /var/www/html/cronjobs/starter2.php
